# systemd/rover.service
[Unit]
Description=Rover Discord Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/rover/rover
ExecStart=/home/rover/rover/.venv/bin/python -m bot.discord_bot
Restart=always
User=rover
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
