# systemd/rover.service
[Unit]
Description=Rover Discord Bot + GPS Logger
After=network.target

[Service]
Type=simple
User=rover
WorkingDirectory=/home/rover/rover
ExecStart=/home/rover/rover/.venv/bin/python /home/rover/rover/main.py
Restart=always
RestartSec=10

# Log en UTF-8
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONIOENCODING=utf-8

[Install]
WantedBy=multi-user.target
