# systemd/rover.service
[Unit]
Description=Rover Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=rover
WorkingDirectory=/home/rover/rover

# Lance directement le bot avec le venv
ExecStart=/home/rover/rover/.venv/bin/python /home/rover/rover/main.py

# Redémarrage automatique en cas de crash
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
